# infra/systemd/mega-quixai-agent-seduction.service
[Unit]
Description=MEGA QUIXAI - Agent Seduction (LangGraph)
After=docker.service
Requires=docker.service
PartOf=mega-quixai.target

[Service]
Type=simple
User=quixai
Group=quixai
WorkingDirectory=/opt/mega-quixai
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=info"
Environment="POSTGRES_HOST=postgres"
Environment="REDIS_HOST=redis"
ExecStart=/opt/mega-quixai/venv/bin/python -m agents.seduction.main
ExecReload=/bin/kill -HUP $MAINPID

# Restart behavior
Restart=on-failure
RestartSec=5s
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryLimit=1G
CPUQuota=50%
CPUAccounting=true

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadOnlyDirectories=/etc:/usr
ProtectHome=yes
ProtectSystem=strict

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mega-quixai-seduction

[Install]
WantedBy=mega-quixai.target
